# build-pkgs-conf/tpls/debian-8/nxs-chat-srv.postinst
#!/bin/sh -e

. /usr/share/debconf/confmodule
db_version 2.0 || [ $? -lt 30 ]

if [ "$1" = "configure" ]; then

	chmod 750 /etc/nxs-chat-srv
	chmod 750 /etc/nxs-chat-srv/ssl
	chmod 640 /etc/nxs-chat-srv/nxs-chat-srv.conf
	chmod 750 /var/log/nxs-chat-srv
	chmod -R 750 /var/spool/nxs-chat-srv

	systemctl enable nxs-chat-srv.service > /dev/null
fi

db_stop

exit 0
